# templates/default/upstart.conf.erb

description     "Google Auth Proxy instance"
author          "Kevin Lamontagne <kevin@demarque.com>"

start on runlevel [2345]
stop on shutdown

respawn

console log

setuid <%= @user %>
exec /opt/go/bin/google_auth_proxy -config=/etc/google_auth_proxy/<%= @service_name %>.conf
